# crontab/每分钟.sh
#!/usr/bin/env bash
* * * * * bash /path/to/file.sh arg1 arg2
